# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file

[supervisord]
logfile=/tmp/supervisord.log ; supervisord log file
logfile_maxbytes=50MB        ; maximum size of logfile before rotation
logfile_backups=10           ; number of backed up logfiles
loglevel=info                ; info, debug, warn, trace
pidfile=/tmp/supervisord.pid ; pidfile location
nodaemon=false               ; run supervisord as a daemon
minfds=1024                  ; number of startup file descriptors
minprocs=200                 ; number of process descriptors

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:reddit_clone]
command=/home/movies_django_rest_api/venv/bin/gunicorn reddit_clone.wsgi:application --bind 0.0.0.0:8000
directory=/home/movies_django_rest_api
autostart=true
autorestart=true
stderr_logfile=/home/movies_django_rest_api/logs/reddit_clone.err.log
stdout_logfile=/home/movies_django_rest_api/logs/reddit_clone.out.log
environment=DJANGO_SETTINGS_MODULE=reddit_clone.settings

[group:reddit_clone]
programs=reddit_clone